# deploy/backend.service
[Unit]
Description=CardioLearn Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/MLTCC/backend
Environment="PATH=/home/ubuntu/MLTCC/backend/venv/bin:/usr/bin"
ExecStart=/home/ubuntu/MLTCC/backend/venv/bin/python -m uvicorn src.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
